# machines/desktop/home/.config/autostart/StreamController.desktop
[Desktop Entry]
Type=Application
Name=StreamController
Exec=flatpak run com.core447.StreamController -b
